# supervisors/laravelEcho.conf
[program:laravelEcho]
command=/usr/bin/laravel-echo-server start --dir=/home/v6market/domains/v6.exchange-demo.ir/webazin
numprocs=1
autostart=true
autorestart=true
user=root

